# Trackermon.service
[Unit]
Description=Trackermon daemon
After= graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/daemon 
User=root

[Install]
WantedBy=graphical.target
